# files/rc.local
cp /etc/resolv.conf.final /etc/resolv.conf
